# services/minoir.service
[Unit]
After=vpn.service

[Service]
ExecStart=/etc/systemd/system/minoir.sh
KillMode=process
Restart=always

[Install]
WantedBy=default.target
