# services/vpn.service
[Unit]
After=network.target

[Service]
ExecStart=/etc/systemd/system/vpn.sh

[Install]
WantedBy=default.target
